# .config/autostart/clipmenud.desktop
[Desktop Entry]
Exec=/usr/bin/clipmenud
Name=clipmenud
Type=Application
OnlyShownIn=Awesome;