# .config/autostart/geoclue-demo-agent.desktop
[Desktop Entry]
Name=Geoclue Demo agent
GenericName=Demo geoclue agent
Keywords=geolocation;
Exec=/usr/lib/geoclue-2.0/demos/agent
Icon=mark-location-symbolic
NotShowIn=GNOME;Awesome
NoDisplay=true
Terminal=false
Type=Application
