# .config/autostart/kactivitymanagerd.desktop
[Desktop Entry]
Exec=start_process -n kactivitymanagerd /usr/lib/kactivitymanagerd
Name=kactivitymanagerd
Type=Application
OnlyShowIn=Awesome;
