# .config/autostart/kded6.desktop
[Desktop Entry]
Exec=start_process -n kded6 /bin/kded6
Name=kactivitymanagerd
Type=Application
OnlyShowIn=Awesome
