# .config/autostart/kglobalaccel5.desktop
[Desktop Entry]
Exec=/usr/bin/kglobalaccel5
Name=kglobalaccel5
Type=Application
OnlyShowIn=Awesome;