# .config/autostart/ksnip.desktop
[Desktop Entry]
Name=ksnip
Icon=ksnip
Exec=ksnip
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
