# .config/autostart/obsidian.desktop
[Desktop Entry]
Name=Obsidian
Exec=/usr/bin/obsidian
Terminal=false
Type=Application
Icon=obsidian
