# .config/autostart/synology-drive.desktop
[Desktop Entry]
Exec=start_process -n cloud-drive-ui synology-drive start
Name=synology-drive
Type=Application
OnlyShownIn=Awesome;
