# .config/autostart/volctl.desktop
[Desktop Entry]
Exec=start_process volctl
Name=volctl
Type=Application
OnlyShowIn=Awesome