# .config/autostart/volctl.desktop
[Desktop Entry]
Exec=sh -c "sleep 5 && volctl"
Name=volctl
Type=Application
OnlyShowIn=Awesome
